# docker/deploy/etc/s6-overlay/s6-rc.d/queue-worker/run
#!/command/with-contenv bash
s6-setuidgid webuser php $WEBUSER_HOME/artisan queue:work --tries=3 --no-ansi -q