# docker/deploy/etc/s6-overlay/s6-rc.d/queue-worker/run
#!/command/with-contenv bash
echo "💪  Starting the queue worker..."
s6-setuidgid webuser php $WEBUSER_HOME/artisan queue:work --tries=3 --no-ansi -q