# docker/deploy/etc/s6-overlay/s6-rc.d/queue-worker/run
#!/command/execlineb -P
s6-setuidgid webuser php /var/www/html/artisan queue:work --tries=3 --no-ansi -q
