# docker/deploy/etc/services.d/queue-worker/run
#!/usr/bin/with-contenv bash

# Switch to web user
su webuser

# Bring up PHP queue worker
exec php /var/www/html/artisan queue:work --tries=3 --no-ansi -q
