# docker/deploy/etc/s6-overlay/s6-rc.d/cron/run
#!/command/execlineb -P
echo "⏰  Starting the cron service..."
cron -f -l 8