# docker/deploy/etc/s6-overlay/s6-rc.d/cron/run
#!/command/with-contenv bash
echo "⏰  Starting the cron service..."
s6-setuidgid root /usr/sbin/cron -f -l 8